# supervisor_production.conf
[program:krv-je-zivot]
command=docker-compose -f production.yml up
directory=/root/krv-je-zivot
redirect_stderr=true
autostart=true
autorestart=true
priority=10
